# launchd/cmux-memory-watchdog/launchd/com.golems.cmux-memory-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.golems.cmux-memory-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>$HOME/Gits/ChatGPTMCPcmux/launchd/cmux-memory-watchdog/bin/cmux-memory-watchdog.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CMUX_MEM_WATCHDOG_FOOTPRINT_THRESHOLD_GB</key>
    <string>5</string>
    <key>CMUX_MEM_WATCHDOG_COMPRESSOR_THRESHOLD_GB</key>
    <string>12</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/tmp/cmux-watchdog/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/cmux-watchdog/launchd.stderr.log</string>
</dict>
</plist>
