# launchd/cmux-contract-nightly/launchd/com.golems.cmux-contract-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.golems.cmux-contract-nightly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/etanheyman/Gits/cmuxlayer/launchd/cmux-contract-nightly/bin/cmux-contract-nightly.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/etanheyman/.local/state/cmux/cmux-contract-nightly.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/etanheyman/.local/state/cmux/cmux-contract-nightly.stderr.log</string>
</dict>
</plist>
